# examples/home-control-housekeeping.timer
[Unit]
Description=Run home-control housekeeping regularly

[Timer]
OnCalendar=*-*-* 0/6:30:00

[Install]
WantedBy=timers.target